# docker/workspace/crontab/laradock
* * * * * laradock php /var/www/moodle/public/admin/cli/cron.php >> /dev/null 2>&1
