# docker/workspace/crontab/laradock.save
* * * * * laradock php /var/www/dule:run >> /dev/null 2>&1
